# profile.d/32-posh-git.ps1

# profile.d/32-posh-git.ps1
# Deprecated: use 32-git-completion.ps1 with $global:PPP.EnableGitModule.
